# deploy/systemd/asset-manager-backup.service
[Unit]
Description=Asset Manager — pg_dump diario
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=rodrigo
WorkingDirectory=/home/rodrigo/asset-manager
ExecStart=/home/rodrigo/asset-manager/deploy/scripts/backup.sh
