# deploy/systemd/benchmarking-backup.service
[Unit]
Description=Benchmarking — pg_dump diario
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=rcoloma
WorkingDirectory=/home/rcoloma/Benchmarking
ExecStart=/home/rcoloma/Benchmarking/deploy/scripts/backup.sh
